# pragmatic.programmer.tips.cli.crontab
30 04 * * * /bin/pragmatic.programmer.tips/pragmatic.programmer.tips.cli
