# pragmatic.programmer.tips.cli.crontab
30 04 * * * dotnet pragmatic.programmer.tips.cli.dll
